# orca-hub.service
[Unit]
Description=OrcaHub (Phoenix dev server - agent mode)
After=network.target postgresql.service

[Service]
Type=simple
User=zach
WorkingDirectory=/home/zach/orca_hub
EnvironmentFile=/home/zach/orca_hub/.env
ExecStart=/usr/bin/elixir --erl "-kernel connect_all false" --name debian_agent@192.168.1.177 --cookie ${RELEASE_COOKIE} -S mix phx.server
Restart=on-failure
RestartSec=5

# Give the BEAM VM time to shut down gracefully
TimeoutStopSec=30

# Logging — stdout/stderr go to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orca-hub

[Install]
WantedBy=multi-user.target
