# systemd/orca-hub.service.template
[Unit]
Description=OrcaHub (compiled OTP release)
After=network-online.target
Wants=network-online.target
StartLimitBurst=3
StartLimitIntervalSec=60

[Service]
Type=simple
User={{USER}}
WorkingDirectory={{HOME}}/orca-hub-releases/current
EnvironmentFile={{HOME}}/orca-hub-releases/.env
Environment="PATH={{HOME}}/.local/share/mise/shims:{{HOME}}/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PHX_SERVER=true"
ExecStart={{HOME}}/orca-hub-releases/current/bin/orca_hub start
Restart=on-failure
RestartSec=5

# Give the BEAM VM time to shut down gracefully
TimeoutStopSec=30

# Logging — stdout/stderr go to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=orca-hub

# Memory limits — BEAM peaks observed around 2.2G
MemoryHigh=2500M
MemoryMax=3G

# Basic hardening. Further hardening (ProtectSystem=strict, ProtectHome=read-only,
# ReadWritePaths=...) is deferred — the app writes to ~/.claude, project
# directories, and other paths that need careful enumeration first.
NoNewPrivileges=yes
PrivateTmp=yes
PrivateDevices=yes

[Install]
WantedBy=multi-user.target
